# app-metrics/starlink_exporter/files/starlink_exporter.initd
#!/sbin/openrc-run
# Copyright 2023 Gentoo authors
# Distributed under the terms of the GNU General Public License v2

description="Prometheus starlink metrics exporter"

command="/usr/bin/starlink_exporter"
command_background="true"
error_log="/var/log/starlink_exporter.log"
output_log="/var/log/starlink_exporter.log"
pidfile=${pidfile:-"/run/starlink_exporter.pid"}

depend() {
	after net
}
